# slm-server/ansible/roles/vnc/templates/vnc-websockify.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
#
# Websockify service for noVNC - managed by Ansible

[Unit]
Description=Websockify for noVNC (AutoBot)
After=network.target {{ vnc_service_name }}.service
Requires={{ vnc_service_name }}.service

[Service]
Type=simple
User={{ vnc_user }}
Group={{ vnc_group }}

# Calculate VNC port from display number (5900 + display)
ExecStart=/usr/bin/websockify --web={{ novnc_path }} 0.0.0.0:{{ websockify_port }} localhost:{{ 5900 + vnc_display | int }}

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
