# scripts/infrastructure/browser-playwright.service
[Unit]
Description=Playwright API Server for Browser VM (Headed Mode)
After=browser-vnc.service
Requires=browser-vnc.service

[Service]
Type=simple
User=autobot
WorkingDirectory=/home/autobot
Environment="HEADLESS=false"
Environment="DISPLAY=:1"
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node /home/autobot/playwright-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
