# autobot-slm-backend/ansible/roles/npu-worker/templates/autobot-npu-worker.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
[Unit]
Description=AutoBot NPU Acceleration Worker
After=network.target

[Service]
Type=simple
User={{ npu_user }}
Group={{ npu_group }}
WorkingDirectory={{ npu_install_dir }}
Environment="PATH={{ npu_install_dir }}/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/autobot/autobot-npu-worker/.env
ExecStart={{ npu_install_dir }}/venv/bin/python {{ npu_install_dir }}/npu-worker.py
Restart=always
RestartSec=10
StandardOutput=append:{{ npu_log_dir }}/npu-worker.log
StandardError=append:{{ npu_log_dir }}/npu-worker-error.log

[Install]
WantedBy=multi-user.target
