# ansible/roles/npu_worker/templates/npu-worker.service.j2
[Unit]
Description=AutoBot NPU Worker - Intel OpenVINO Hardware Acceleration
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=autobot
Group=autobot
WorkingDirectory=/home/autobot/backend

# Environment
Environment="OPENVINO_DEVICE={{ openvino_device }}"
Environment="NPU_WORKER_PORT={{ npu_worker_port }}"
Environment="WORKER_THREADS={{ worker_threads }}"
Environment="INFERENCE_OPTIMIZATION={{ inference_optimization }}"

# Start command (placeholder - will be updated when code is synced)
ExecStart=/usr/bin/python3 -m uvicorn npu_worker:app --host 0.0.0.0 --port {{ npu_worker_port }}

# Restart policy
Restart=always
RestartSec=5s

# Limits
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobot-npu-worker

[Install]
WantedBy=multi-user.target
