# autobot-slm-backend/ansible/roles/browser/templates/playwright.service.j2
[Unit]
Description=AutoBot Playwright Browser Automation
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=autobot
Group=autobot
WorkingDirectory=/opt/autobot/autobot-browser-worker

# Environment
Environment="PLAYWRIGHT_PORT={{ playwright_port }}"
Environment="BROWSER_TIMEOUT={{ browser_timeout }}"
Environment="HEADLESS_MODE={{ headless_mode }}"
Environment="DISPLAY={{ vnc_display }}"

# Start command - Node.js Express server
ExecStart=/usr/bin/node /opt/autobot/autobot-browser-worker/playwright-server.js

# Restart policy
Restart=always
RestartSec=5s

# Limits
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobot-playwright

[Install]
WantedBy=multi-user.target
