# autobot-slm-backend/ansible/roles/browser/templates/playwright.service.j2
[Unit]
Description=AutoBot Playwright Browser Automation
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=autobot
Group=autobot
WorkingDirectory=/home/autobot/backend

# Environment
Environment="PLAYWRIGHT_PORT={{ playwright_port }}"
Environment="BROWSER_TIMEOUT={{ browser_timeout }}"
Environment="HEADLESS_MODE={{ headless_mode }}"
Environment="DISPLAY={{ vnc_display }}"

# Start command (placeholder - will be updated when code is synced)
ExecStart=/usr/bin/python3 -m uvicorn playwright_api:app --host 0.0.0.0 --port {{ playwright_port }}

# Restart policy
Restart=always
RestartSec=5s

# Limits
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobot-playwright

[Install]
WantedBy=multi-user.target
