# slm-server/ansible/roles/browser-automation/templates/browser-automation.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
[Unit]
Description=AutoBot Browser Automation (Playwright)
After=network.target

[Service]
Type=simple
User={{ browser_user }}
Group={{ browser_group }}
WorkingDirectory={{ browser_install_dir }}
Environment="DISPLAY=:99"
ExecStartPre=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 &
ExecStart=/usr/bin/npx playwright run-server --host {{ browser_host }} --port {{ browser_port }}
Restart=always
RestartSec=10
StandardOutput=append:{{ browser_log_dir }}/browser.log
StandardError=append:{{ browser_log_dir }}/browser-error.log

[Install]
WantedBy=multi-user.target
