# slm-server/ansible/roles/monitoring/templates/prometheus.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
#
# Prometheus Systemd Service - Generated by Ansible

[Unit]
Description=Prometheus Monitoring System
Documentation=https://prometheus.io/docs/introduction/overview/
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_group }}
ExecReload=/bin/kill -HUP $MAINPID
ExecStart={{ prometheus_install_dir }}/prometheus \
  --config.file={{ prometheus_config_dir }}/prometheus.yml \
  --storage.tsdb.path={{ prometheus_data_dir }} \
  --storage.tsdb.retention.time={{ prometheus_retention_time }} \
  --storage.tsdb.retention.size={{ prometheus_retention_size }} \
  --web.console.templates={{ prometheus_install_dir }}/consoles \
  --web.console.libraries={{ prometheus_install_dir }}/console_libraries \
  --web.listen-address=0.0.0.0:{{ prometheus_port }} \
  --web.enable-lifecycle

SyslogIdentifier=prometheus
Restart=always
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ prometheus_data_dir }}

[Install]
WantedBy=multi-user.target
