# autobot-slm-backend/ansible/roles/distributed_setup/templates/autobot-coordinator.service.j2
# AutoBot Coordinator Service
# Generated by Ansible - DO NOT EDIT MANUALLY

[Unit]
Description=AutoBot Distributed Coordinator
After=network.target redis.target
Wants=autobot-health-monitor.timer

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ project_root }}
EnvironmentFile={{ config_dir }}/distributed/environment.env

ExecStart={{ project_root }}/venv/bin/python -m autobot_user_backend.main

Restart={{ restart_policy }}
RestartSec={{ restart_sec }}
StartLimitInterval=5min
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths={{ logs_dir }} {{ backup_dir }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ syslog_identifier }}-coordinator

[Install]
WantedBy=multi-user.target
