# autobot-slm-backend/ansible/roles/backend/templates/autobot-backend.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
[Unit]
Description=AutoBot Backend API
After=network.target redis.service

[Service]
Type=simple
User={{ backend_user }}
Group={{ backend_group }}
WorkingDirectory={{ backend_install_dir }}
Environment="PATH={{ backend_install_dir }}/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile={{ backend_install_dir }}/.env
EnvironmentFile=-{{ service_auth_keys_dir | default('/etc/autobot/service-keys') }}/{{ service_auth_service_id | default('main-backend') }}.env
ExecStart={{ backend_install_dir }}/venv/bin/uvicorn main:app --host {{ backend_host }} --port {{ backend_port }} --workers {{ backend_workers }}
Restart=always
RestartSec=5
StandardOutput=append:{{ backend_log_dir }}/backend.log
StandardError=append:{{ backend_log_dir }}/backend-error.log

[Install]
WantedBy=multi-user.target
