# autobot-slm-backend/ansible/roles/centralized_logging/templates/loki.service.j2
[Unit]
Description=Loki Log Aggregation System
Documentation=https://grafana.com/docs/loki/latest/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ logging_user }}
Group={{ logging_group }}
ExecStart={{ loki_binary_path }} -config.file={{ loki_config_dir }}/config.yml
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loki

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths={{ loki_data_dir }}

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
