# autobot-slm-backend/ansible/roles/agent_config/templates/autobot-agent.service.j2
# AutoBot Agent Service
# Generated by Ansible - DO NOT EDIT MANUALLY

[Unit]
Description=AutoBot Agent ({{ agent_role }})
After=network.target redis.target postgresql.target
{% if install_ollama %}
Wants=ollama.service
{% endif %}

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ project_root }}
EnvironmentFile={{ config_dir }}/agent.env

ExecStart={{ venv_dir }}/bin/python -m autobot_agent.main

Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths={{ logs_dir }} {{ uploads_dir }} {{ models_dir }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobot-agent-{{ agent_role }}

[Install]
WantedBy=multi-user.target
