# slm-server/ansible/roles/slm-agent/templates/slm-agent.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
[Unit]
Description=SLM Fleet Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ slm_agent_user }}
Group={{ slm_agent_group }}
WorkingDirectory=/opt/slm-agent
ExecStart=/usr/bin/python3 /opt/slm-agent/slm-agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
