# ansible/roles/ai_stack/templates/ai-stack.service.j2
[Unit]
Description=AutoBot AI Stack - AI Processing Services
After=network.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=autobot
Group=autobot
WorkingDirectory=/home/autobot/backend

# Environment
Environment="AI_STACK_PORT={{ ai_stack_port }}"
Environment="OLLAMA_HOST=http://127.0.0.1:{{ ollama_port }}"
Environment="MODEL_CACHE_SIZE={{ model_cache_size }}"
Environment="INFERENCE_THREADS={{ inference_threads }}"

# Start command (placeholder - will be updated when code is synced)
ExecStart=/usr/bin/python3 -m uvicorn ai_stack:app --host 0.0.0.0 --port {{ ai_stack_port }}

# Restart policy
Restart=always
RestartSec=5s

# Limits
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobot-ai-stack

[Install]
WantedBy=multi-user.target
