# autobot-slm-backend/ansible/roles/slm_agent/templates/slm-agent.service.j2
# AutoBot - AI-Powered Automation Platform
# Copyright (c) 2025 mrveiss
# Author: mrveiss
#
# SLM Node Agent systemd service unit
# Managed by Ansible - do not edit manually
# Node: {{ slm_node_id }}
# Generated: {{ ansible_date_time.iso8601 }}
#
# Auto-detect if agent runs on same host as SLM backend - use localhost if so
{% set admin_host = slm_admin_url | regex_replace('^https?://([^:/]+).*$', '\\1') %}
{% if admin_host == ansible_host or admin_host == inventory_hostname %}
{% set effective_admin_url = slm_admin_url | regex_replace(admin_host, '127.0.0.1') %}
{% else %}
{% set effective_admin_url = slm_admin_url %}
{% endif %}

[Unit]
Description=SLM Agent - Service Lifecycle Manager
Documentation=https://github.com/mrveiss/AutoBot-AI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ slm_agent_user }}
Group={{ slm_agent_group }}
WorkingDirectory={{ slm_agent_dir }}

# Environment variables
Environment=SLM_NODE_ID={{ slm_node_id }}
Environment=SLM_ADMIN_URL={{ effective_admin_url }}
Environment=PYTHONPATH={{ slm_agent_dir }}
Environment=PYTHONUNBUFFERED=1

# Agent execution
ExecStart=/usr/bin/python3 -m slm.agent.agent \
    --admin-url {{ effective_admin_url }} \
    --node-id {{ slm_node_id }} \
    --interval {{ slm_heartbeat_interval }}{% if slm_services_to_monitor | length > 0 %} \
    --services {{ slm_services_to_monitor | join(' ') }}{% endif %}{% if slm_agent_debug %} \
    --debug{% endif %}


# Restart policy
Restart=always
RestartSec={{ slm_agent_restart_sec }}

# Watchdog (systemd will restart if agent becomes unresponsive)
WatchdogSec={{ slm_agent_watchdog_sec }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ slm_agent_data_dir }}
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
